# supervisord.conf
[supervisord]
nodaemon=true

[group:clairjwt]
programs=clair,jwtproxy

[program:jwtproxy]
command=/usr/bin/jwtproxy -config /config/config.yaml
autorestart=true
startretries=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:clair]
command=clair -config /config/config.yaml
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
